# scripts/com.homelab-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.homelab-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BINARY_PATH__</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MQTT_BROKER</key>
        <string>gigantic.lan</string>
        <key>TOPIC_PREFIX</key>
        <string>agent</string>
        <key>SERVICES</key>
        <string>__SERVICES__</string>
        <key>DEPLOY_DIR</key>
        <string>__DEPLOY_DIR__</string>
<key>METRICS_INTERVAL_SECONDS</key>
        <string>60</string>
        <key>DEBUG</key>
        <string>false</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/homelab-agent.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/homelab-agent.log</string>
</dict>
</plist>
